# templates/jenkins-agent.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-/Apple/DTD PLIST 1.0/EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.osrfoundation.build.swarm-client</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/java</string>
		<string>-jar</string>
		<string><%= @swarm_jar_path %></string>
		<string>-url</url>
		<string><%= @jenkins_url %></string>
		<string>-name</url>
		<string><%= @agent_name %></string>
		<string>-username</string>
		<string><%= @username %></string>
		<string>-password</string>
		<string><%= @token %></string>
		<string>-description</string>
		<string><%= @description %></string>
		<string>-mode</string>
		<string>exclusive</string>
		<string>-executors</url>
		<string>1</string>
		<string>-fsroot</string>
		<string>/Users/jenkins</string>
		<string>-labels</url>
		<string><%= @labels.join(" ") %></string>
		<string>-disableClientsUniqueId</string>
		<string>-deleteExistingClients</string>
	</array>
	<string>org.osrfoundation.build.jenkins-agent</string>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>jenkins</string>
	<key>WorkingDirectory</key>
	<string>/Users/Jenkins</string>
	<key>StandardInPath</key>
	<string>/dev/null</string>
	<key>StandardErrorPath</key>
	<string>jenkins-agent.err.log</string>
	<key>StandardOutPath</key>
	<string>jenkins-agent.out.log</string>
	<true/>
</dict>
